# comms/lcdproc/pkg/lcdproc.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/bin/lcdproc"

. /etc/rc.d/rc.subr

rc_cmd $1
